# pkg/debian/salt-master.postinst
chown -R salt:salt /etc/salt /var/log/salt /opt/saltstack/salt/ /var/cache/salt/ /var/run/salt
systemctl enable salt-master
