# pkg/common/salt-api.upstart
description "Salt API"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

script
exec salt-api
end script
