# pkg/old/openbsd/salt-syncdic.rc-d
#!/bin/sh

daemon="/usr/local/bin/salt-syncdic"
daemon_flags="-d"

. /etc/rc.d/rc.subr

pexp="/usr/local/bin/python2.7 $daemon"

rc_cmd $1
