# tests/integration/files/file/base/supervisor.conf
[unix_http_server]
file={{ supervisor_sock }}       ; (the path to the socket file)
chmod=0700                       ; socket file mode (default 0700)

[supervisord]
logfile={{ virtual_env }}/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile={{ virtual_env }}/supervisord.pid ; (supervisord pidfile;default supervisord.pid)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://{{ supervisor_sock }} ; use a unix:// URL  for a unix socket

[program:sleep_service]
command=sleep 600
{% if not autostart %}
autostart=false
{% endif %}

[program:sleep_service2]
command=sleep 600
{% if not autostart %}
autostart=false
{% endif %}
