# debian/salt-syndic.upstart
description "salt-syndic"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]


exec /usr/bin/salt-syndic >/dev/null 2>&1
