# pkg/salt-syndic.upstart
description "salt-syndic"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

exec su -c salt-syndic
