# pkg/salt-master.upstart
description "Salt Master"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

kill timeout 10

exec /usr/bin/salt-master >/dev/null 2>&1
